# target/linux/mediatek/base-files/etc/hotplug.d/net/30-ext-devices-rps-setting
#!/bin/sh

NPROCS="$(grep -c "^processor.*:" /proc/cpuinfo)"
[ "$NPROCS" -gt 1 ] || exit

IFNAME_PREFIX="${INTERFACE%%[0-9]*}"

if [ "${IFNAME_PREFIX}" = "usb" ] || [ "${IFNAME_PREFIX}" = "wwan" ]; then
    if [ "$ACTION" = add ]; then
        # set usb/wwan rps to cpu0
        echo 1 > /sys/class/net/${INTERFACE}/queues/rx-0/rps_cpus
    fi
    [ -f /etc/init.d/mtk_smp ] && /etc/init.d/mtk_smp restart
fi
